# catalogue/templates/catalogue.service
[Unit]
Description = Catalogue Service

[Service]
User=roboshop
Environment=MONGO=true
<<<<<<< HEAD
Environment=MONGO_URL="mongodb://mongodb.{{ env }}.meppk.xyz:27017/catalogue"
=======
Environment=MONGO_URL="mongodb://mongodb.dev.meppk.xyz:27017/catalogue"
>>>>>>> origin/main
ExecStart=/bin/node /app/server.js
SyslogIdentifier=catalogue

[Install]
WantedBy=multi-user.target